# system/air-mq7.service
[Unit]
Description=Air

[Service]
User=pi
ExecStart=/home/pi/air-quality/mq7.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
