# system/air-mq.service
[Unit]
Description=Air

[Service]
User=dpavlin
ExecStart=/home/dpavlin/air-quality/mq.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
