# system/air-zph02.service
[Unit]
Description=Air

[Service]
User=pi
ExecStart=/home/pi/air-quality/zph02.pl
Restart=always
RestartSec=3s

[Install]
WantedBy=multi-user.target
