# system/air-mq7.service
[Unit]
Description=Air

[Service]
User=pi
ExecStart=/home/pi/air-quality/mq7.sh
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=network-online.target
